# system/systemd/winewarden.socket
[Unit]
Description=WineWarden Socket

[Socket]
ListenStream=%t/winewarden/winewarden.sock

[Install]
WantedBy=sockets.target
